# systemd/aqua-matrix-heartbeat.service
[Unit]
Description=Aqua Matrix heartbeat (10-minute status DM to default target)
Documentation=https://github.com/inblockio/aqua-matrix-hello
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# siwx-oidc mints a fresh device_id on every re-auth, but the matrix-sdk
# SQLite crypto store binds to the previous one — so any restart fails
# with `account in the store doesn't match`. Wiping the crypto SQLite files
# (NOT config.toml) before each start sidesteps this so /restart works.
ExecStartPre=/bin/sh -c 'rm -f %h/.aqua-matrix-heartbeat/matrix-sdk-*.sqlite3*'
ExecStart=%h/aqua-matrix-hello/target/debug/aqua-matrix-agent --heartbeat --key-file heartbeat.pem --store-dir %h/.aqua-matrix-heartbeat
Restart=on-failure
RestartSec=30s
# Run in the repo dir so the default agent.pem key file is picked up
WorkingDirectory=%h/aqua-matrix-hello
# Forward CONTEXT_WINDOW for accurate context% reporting; matches settings.json
Environment=CONTEXT_WINDOW=1000000
# Log to journal under the unit name
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
