# systemd/aqua-matrix-heartbeat.service
[Unit]
Description=Aqua Matrix heartbeat (10-minute status DM to default target)
Documentation=https://github.com/inblockio/aqua-matrix-hello
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Device-id stability is handled in code via session caching in
# AgentClient::connect (see docs/ARCHITECTURE.md). When the cached access
# token is still valid (~5 min window), restart reuses the same device_id
# and the SQLite crypto store opens cleanly. When it has expired, fresh
# auth gets a new device_id and the in-code store-mismatch handler wipes
# matrix-sdk-*.sqlite3* and retries restore_session once.
ExecStart=%h/aqua-matrix-hello/target/debug/aqua-matrix-agent --heartbeat --key-file heartbeat.pem --store-dir %h/.aqua-matrix-heartbeat
Restart=on-failure
RestartSec=30s
# Run in the repo dir so the default agent.pem key file is picked up
WorkingDirectory=%h/aqua-matrix-hello
# Forward CONTEXT_WINDOW for accurate context% reporting; matches settings.json
Environment=CONTEXT_WINDOW=1000000
# Log to journal under the unit name
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
