# systemd/aqua-matrix-heartbeat.service
[Unit]
Description=Aqua Matrix heartbeat (10-minute status DM to default target)
Documentation=https://github.com/inblockio/aqua-matrix-agent
After=network-online.target
Wants=network-online.target
# Crash-loop guard belongs in [Unit] in modern systemd; [Service] placement
# is silently ignored.
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
# Device-id stability is handled in code via session caching in
# AgentClient::connect (see docs/ARCHITECTURE.md). The daemon also rotates
# its matrix-sdk Client in-process ahead of every token expiry, so steady
# state should no longer need a systemd restart to recover from
# M_UNKNOWN_TOKEN.
# Dedicated ops binary (was `aqua-matrix-agent --heartbeat` before the
# workspace split; the relay lifecycle now lives in aqua-matrix-relay).
ExecStart=%h/aqua-matrix-agent/target/debug/aqua-matrix-heartbeat --key-file heartbeat.pem --store-dir %h/.aqua-matrix-heartbeat
Restart=always
RestartSec=5s
# Run in the repo dir so the default agent.pem key file is picked up
WorkingDirectory=%h/aqua-matrix-agent
# Forward CONTEXT_WINDOW for accurate context% reporting; matches settings.json
Environment=CONTEXT_WINDOW=1000000
# Per-instance config (AGENT_TARGET et al.): loaded from this machine-local
# `.env` at startup. Not in the repo — copy `.env.example` and set your target.
Environment=AGENT_ENV_FILE=%h/.aqua-matrix-heartbeat/agent.env
# Log to journal under the unit name
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
